# etc/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/fd/1
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid 

[program:gunicorn]
command=gunicorn --ssl-version 2 --log-config /etc/gunicorn_logging.conf --config spectre/gunicorn_config.py spectre.wsgi
user=gunicorn
directory=/spectre
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
